# com.ragmonitor.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ragmonitor.startup</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mhanheide/workspace/qdrant_file_scanner/start_rag_monitor.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/mhanheide/Library/Logs/rag-monitor-launchd-out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mhanheide/Library/Logs/rag-monitor-launchd-err.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/mhanheide/workspace/qdrant_file_scanner</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
